# docker/supervisor/supervisor.conf
[unix_http_server]
file=/run/supervisord.sock

# Cron via supervisor
# https://stackoverflow.com/a/29565906/1756667
# https://gist.github.com/martinrusev/7015e393d46647dbad15
[program:cron]
command=cron -f -L 15
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=10
stopwaitsecs=5

# Supervisord configuration
# http://supervisord.org/logging.html
[supervisord]
nodaemon=true
logfile=/app/data/files/private/logs/supervisor/supervisord.log
childlogdir=/app/data/files/private/logs/supervisor
logfile_maxbytes=20MB
logfile_backups=5

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock
